# supervisor-app.conf
[program:npm-build]
command = /var/www/html/install.sh
startsecs = 0
autorestart = false
startretries = 1

[program:nginx-app]
command = /usr/sbin/nginx
