# supervisor-app.conf
[program:npm-build]
command = /bin/sh -c "pwd && npm run build"
directory = /var/www/html

[program:nginx-app]
command = /usr/sbin/nginx
